# systemd/almalex-sync-fedlex.service
[Unit]
Description=Sync Fedlex documents for almalex.ch
After=network.target

[Service]
Type=oneshot
User=joshua
Group=joshua
WorkingDirectory=/var/www/almalex.ch/backend
Environment=PATH=/var/www/almalex.ch/backend/.venv/bin
ExecStart=/var/www/almalex.ch/backend/.venv/bin/python -m cli.main sync-fedlex --embed

[Install]
WantedBy=multi-user.target
